# client/achillebackups.service
[Unit]
Description=Executes AchilleBackups script

[Service]
ExecStart=/usr/local/bin/achillebackups automated
